# units/roundtable.service
[Unit]
Description=Roundtable daily session

[Service]
Type=oneshot
ExecStart=%h/.local/bin/roundtable session --with-games
