# units/roundtable-bind.service
[Unit]
Description=Roundtable weekly bind

[Service]
Type=oneshot
ExecStart=%h/.local/bin/roundtable bind
